# etc/profile.d/init_aliases.sh.new
alias dmesg='dmesg --color=auto'
alias grep='grep --color=auto'

alias bc='bc -lq'


alias emacs='emacs -nw'
