# etc/profile.d/init_env.sh.new
VISUAL=$(echo /usr/bin/emacs-*-no-x11)
EDITOR=$VISUAL
PAGER='most'

HISTCONTROL='ignoredups'

export VISUAL EDITOR PAGER HISTCONTROL
