# etc/profile.d/init_env.sh.new
TERM='linux'
VISUAL='emacs -nw'
EDITOR='emacs -nw'
PAGER='most'

HISTCONTROL='ignoredups'

export TERM VISUAL EDITOR PAGER HISTCONTROL
