# etc/profile.d/init_prompt.sh.new
#!/bin/sh

if [ "`id -u`" = "0" ]; then
    PS1="\[\033[01;31m\]\u\[\033[m\]@\h:\w# "
else
    PS1="\[\033[01;32m\]\u\[\033[m\]@\h:\w$ "
fi
PS2='~> '

export PS1 PS2
