# etc/profile.d/init_aliases.sh.new
alias grep='grep --color=auto'

alias bc='bc -lq'

alias emacs='emacs -nw'
