# etc/profile.d/init_env.sh.new
VISUAL='emacs -nw'
EDITOR='emacs -nw'
PAGER='most'

export VISUAL EDITOR PAGER
